# services/mbot-scratch-gui.service
[Unit]
Description=Server for the MBot Scratch GUI
After=network.target

[Service]
ExecStart=sudo npx -y serve -s /data/www/scratch/ -l 8602
[Install]
WantedBy=multi-user.target
